# contrib/chef-cookbooks/cookbooks/cloud_servers_vpc/templates/default/cloud_servers_vpc.cron.erb
*/15 * * * * cloud_servers_vpc <%= node[:cloud_servers_vpc][:home] %>/contrib/bin/cloud-servers-vpc-cleanup
